# docker/supervisor.conf.d/janitoo_datalog_rrd.conf
[program:jnt_datalog_rrd]
user=root
autostart=true
autorestart=true
command=jnt_datalog_rrd -c /opt/etc/janitoo_datalog_rrd.conf
